# systemd/goappnixosmodule.service.in
[Unit]
Description=A goappnixosmodule service

[Service]
DynamicUser=yes
ExecStart=@goappnixosmodule@/bin/goappnixosmodule

[Install]
WantedBy=multi-user.target
